# replication.service
[Unit]
Description=Replication service

[Service]
Type=simple
ExecStart=/opt/bin/replication.sh
